# config/paleblue.systems/a3/etc/systemd/system/polkadot-0.service
[Unit]
Description=polkadot 0 0
Wants=network-online.target
After=network-online.target

[Service]
User=substrate
Group=substrate
Type=simple
ExecStart=/usr/local/bin/polkadot \
  --pruning archive \
  --name 'a3 0' \
  --base-path /var/lib/substrate/polkadot-staging-0 \
  --chain /usr/share/substrate/polkadot-staging.json \
  --force-polkadot \
  --port 30433 \
  --ws-port 7944 \
  --rpc-port 7933 \
  --rpc-cors all \
  --rpc-methods safe \
  --log="parachain=debug" \
  --ws-max-connections 1000 \
  --telemetry-url 'wss://api.telemetry.pelagos.systems/submit/ 0'
Restart=always
RestartSec=120

[Install]
WantedBy=multi-user.target
