# static/etc/systemd/system/nginx.service.d/hourly-restart.conf
[Service]
RuntimeMaxSec=3600
Restart=always
