# config/paleblue.systems/v1/etc/systemd/system/polkadot-1.service
[Unit]
Description=polkadot 1
Wants=network-online.target
After=network-online.target

[Service]
User=root
Group=root
Type=simple
ExecStart=/usr/local/bin/polkadot \
  --validator \
  --name v1-1 \
  --base-path /var/lib/polkadot1 \
  --node-key-file /var/lib/polkadot1/node-key \
  --chain /usr/share/polkadot/relaychain-spec.json \
  --force-polkadot \
  --port 30434 \
  --ws-port 7945 \
  --rpc-port 7934 \
  --rpc-cors all \
  --rpc-methods safe \
  --log="parachain=debug" \
  --ws-max-connections 1000 \
  --telemetry-url 'wss://api.telemetry.pelagos.systems/submit/ 0'
Restart=always
RestartSec=120

[Install]
WantedBy=multi-user.target
