# config/paleblue.systems/v2/etc/systemd/system/polkadot-0.service
[Unit]
Description=polkadot 0
Wants=network-online.target
After=network-online.target

[Service]
User=root
Group=root
Type=simple
ExecStart=/usr/local/bin/polkadot \
  --validator \
  --name v2-0 \
  --base-path /var/lib/polkadot0 \
  --node-key-file /var/lib/polkadot0/node-key \
  --chain /usr/share/polkadot/relaychain-spec.json \
  --force-polkadot \
  --port 30433 \
  --ws-port 7944 \
  --rpc-port 7933 \
  --rpc-cors all \
  --rpc-methods safe \
  --log="parachain=debug" \
  --ws-max-connections 1000 \
  --telemetry-url 'wss://api.telemetry.pelagos.systems/submit/ 0'
Restart=always
RestartSec=120

[Install]
WantedBy=multi-user.target
